# fixtures/clean_b186_plugins_var_ignored/.config/systemd/user/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway

[Service]
ExecStart=/usr/bin/openclaw gateway run
Environment=OPENCLAW_BUNDLED_PLUGINS_DIR=/nonexistent/relocated-plugins

[Install]
WantedBy=default.target
