# fixtures/clean_b150_systemd_restart/.config/systemd/user/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway

[Service]
ExecStart=/usr/bin/openclaw gateway start
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
