# fixtures/clean_b186_no_bundled_override/.config/systemd/user/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/openclaw gateway run
Restart=always
Environment=OPENCLAW_GATEWAY_PORT=8899
Environment=OPENCLAW_SERVICE_KIND=gateway
Environment="OPENCLAW_WINDOWS_TASK_NAME=OpenClaw Gateway"

[Install]
WantedBy=default.target
